# ww.conf
[program:ww]
command=python3 /home/curator2/ww/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ww/test.err.log
stdout_logfile=/home/curator2/ww/test.out.log